# root/etc/services.d/verysync/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /verysync/verysync -home /verysync/config -gui-address "0.0.0.0:$PORT" -logflags=0
